# autostart/onedrive.desktop
[Desktop Entry]
Name=Mount OneDrive
Exec=/home/kyiimn/.local/bin/MountOnedrive.sh
Terminal=false
Type=Application
